# assets/jupyterlab/jupyterlab.service
[Unit]
Description=JupyterLab
After=syslog.target network.target

[Service]
ExecStart=jupyter lab

[Install]
WantedBy=multi-user.target